# root/etc/services.d/dockerd/run
#!/usr/bin/with-contenv bash

# Start docker
exec \
    /usr/bin/dockerd
